# bin/ansible.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/S02ansible

[Install]
WantedBy=default.target